# systemd/forpost-stream-web.service
[Unit]
Description=FORPOST Stream Web Configuration Interface
After=network.target

[Service]
Type=simple
User=rpidrone
WorkingDirectory=__INSTALL_DIR__/web
ExecStart=/usr/bin/python3 __INSTALL_DIR__/web/web_config.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
