# systemd/forpost-stream-restart.service
[Unit]
Description=Restart Forpost Stream on config change

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart forpost-stream.service
